# misc/conf/supervisor.conf
[program:xinhua]
command=python YOUR-PATH/app.py
directory=YOUR-PATH/
autostart=true
